# N1-LEDE/files/etc/crontabs/root
0 4 * * * fstrim -a
16 */6 * * * /usr/bin/cloudflarespeedtest/cloudflarespeedtest.sh start
